# upstart/docker-caddy.conf
description "Caddy container"
start on filesystem and started docker
stop on runlevel [!2345]
respawn
script
  /usr/bin/docker start -a caddy
end script
post-stop script
  /usr/bin/docker stop -t 2 caddy
end script
